# tools/ansible-deprecated/roles/launcher/templates/CentOS/launcher.service.j2
[Unit]
Description=ParaViewWeb Launcher

[Service]
User={{pvw_user}}
Type=simple
ExecStart=/usr/bin/launcher.sh
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
